# playbooks/roles/supervisor/templates/conf.d/mailhog.conf.j2
[program:mailhog]
command=/usr/local/bin/mailhog
autostart=true
autorestart=true
stderr_logfile = {{ supervisor_log_path }}/mailhog-stderr.log