# playbooks/roles/supervisor/templates/conf.d/programs.conf.j2
{% for program in supervisor_programs %}
[program:{{ program.name }}]
{% for i in program if i != 'name' %}
{{ i }}={{ program[i] }}
{% endfor %}
{% endfor %}
