# playbooks/roles/supervisor/templates/conf.d/programs.conf.j2
{% for program in supervisor_programs|sort %}
[program:{{program}}]
{% for key in supervisor_programs[program]|sort %}
{{ key }} = {{ supervisor_programs[program][key] }}
{% endfor %}

{% endfor %}
